# docker/production/etc/s6-overlay/s6-rc.d/init-seeder/up
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv
if { test "${SEEDER_ENABLED}" = "true" }
foreground {
    cd /var/www/html
    php
    artisan
    db:seed
    --force
}
ifnot {
    # If seeder is disabled, just exit gracefully
    exit 0
}


