# docker/production/etc/s6-overlay/s6-rc.d/db-migration/up
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv
if { test "${MIGRATION_ENABLED}" = "true" }
foreground {
    cd /var/www/html
    php
    artisan
    migrate
    --force
}
ifnot {
    # If migration is disabled, just exit gracefully
    exit 0
}

