# docker/production/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/execlineb -P

# Use with-contenv to ensure environment variables are available
with-contenv
if { test "${SCHEDULER_ENABLED}" = "true" }
foreground {
    cd /var/www/html
    php
    artisan
    schedule:work
}
ifnot {
    # If scheduler is disabled, just exit gracefully
    exit 0
} 